# packaging/systemd/zeb-tunnel.service
[Unit]
Description=Zeb public Cloudflare Tunnel
After=network-online.target zeb.service
Wants=network-online.target
Requires=zeb.service

[Service]
Type=simple
User=zeb
Group=zeb
EnvironmentFile=/etc/zeb/zeb.env
Environment=ZEB_ENV_FILE=/etc/zeb/zeb.env
ExecStart=/usr/local/libexec/zeb-tunnel
Restart=always
RestartSec=5
TimeoutStartSec=0
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
