# docker/s6-rc.d/zeb-tunnel/run
#!/command/with-contenv sh
# shellcheck shell=sh
# Cloudflare tunnel for the dashboard.
#
# Two modes:
#   - Named tunnel (ZEB_TUNNEL_TOKEN set): a permanent Cloudflare Tunnel
#     bound to your own domain(s). Public-hostname -> local-service routing
#     lives in the Cloudflare dashboard (Zero Trust > Networks > Tunnels),
#     not in this container, so the URL is stable across rebuilds/restarts.
#   - Quick tunnel (no token): falls back to a free, random
#     https://<random>.trycloudflare.com URL that works from ANY device
#     with no Cloudflare account. The URL changes every restart.
#
# Enabled automatically whenever the dashboard is enabled. Set ZEB_TUNNEL=0
# to disable and rely on the direct host:port bind instead.
#
# The dashboard is still password-protected behind the tunnel (the auth gate
# requires the basic-auth credentials printed by 016-dashboard-auto-auth), so
# a public URL does not mean an open dashboard.

# Only run when the dashboard itself is enabled.
case "${ZEB_DASHBOARD:-}" in
    1|true|TRUE|True|yes|YES|Yes) ;;
    *) exit 0 ;;
esac

# Opt-out switch (default on).
case "${ZEB_TUNNEL:-1}" in
    0|false|FALSE|False|no|NO|No)
        echo "[tunnel] ZEB_TUNNEL disabled — skipping public link." >&2
        exit 0
        ;;
esac

if ! command -v cloudflared >/dev/null 2>&1; then
    echo "[tunnel] cloudflared is not installed in this image; no public link." >&2
    exit 0
fi

port="${ZEB_DASHBOARD_PORT:-9119}"
target="http://localhost:${port}"
user="${ZEB_DASHBOARD_BASIC_AUTH_USERNAME:-admin}"
pw="${ZEB_DASHBOARD_BASIC_AUTH_PASSWORD:-(see DASHBOARD_LOGIN.txt)}"

rule=$(printf '%68s' '' | tr ' ' '=')
line() { printf '|%-68s|\n' "$1"; }

print_link_box() {
    _url="$1"
    echo ""
    echo ""
    printf '+%s+\n' "$rule"
    line ""
    line "    ****  YOUR DASHBOARD LINK  ****"
    line ""
    line "    ------------------------------------------------------------"
    line ""
    line "      OPEN THIS FROM ANY DEVICE:"
    line ""
    line "      $_url"
    line ""
    line "    ------------------------------------------------------------"
    line ""
    line "      USERNAME:   $user"
    line "      PASSWORD:   $pw"
    line ""
    line "    (link may take ~10 seconds to become reachable)"
    line ""
    printf '+%s+\n' "$rule"
    echo ""
    echo ""
}

# Multi-hostname variant for named tunnels, which can front more than one
# custom domain from a single cloudflared process.
print_link_box_multi() {
    echo ""
    echo ""
    printf '+%s+\n' "$rule"
    line ""
    line "    ****  YOUR DASHBOARD LINK(S)  ****"
    line ""
    line "    ------------------------------------------------------------"
    line ""
    for _u in "$@"; do
        line "      $_u"
    done
    line ""
    line "    ------------------------------------------------------------"
    line ""
    line "      USERNAME:   $user"
    line "      PASSWORD:   $pw"
    line ""
    line "    (link may take ~10 seconds to become reachable)"
    line ""
    printf '+%s+\n' "$rule"
    echo ""
    echo ""
}

# ---------- named tunnel (permanent, your own domain) ----------
# ZEB_TUNNEL_TOKEN comes from `cloudflared tunnel token <name>` on any
# machine authenticated to your Cloudflare account. It's a self-contained
# credential — no credentials.json file needs to be mounted into the
# container. Which hostnames route to this tunnel is configured once in
# the Cloudflare dashboard (Zero Trust > Networks > Tunnels > Public
# Hostnames), so it survives container rebuilds untouched.
if [ -n "${ZEB_TUNNEL_TOKEN:-}" ]; then
    echo "[tunnel] starting named Cloudflare tunnel -> ${target}"

    if [ -n "${ZEB_TUNNEL_HOSTNAMES:-}" ]; then
        urls=""
        old_ifs="$IFS"
        IFS=','
        for h in $ZEB_TUNNEL_HOSTNAMES; do
            h_trimmed=$(printf '%s' "$h" | sed 's/^ *//;s/ *$//')
            [ -n "$h_trimmed" ] || continue
            urls="${urls}https://${h_trimmed}
"
        done
        IFS="$old_ifs"
        if [ -n "$urls" ]; then
            printf '%s' "$urls" > /opt/data/DASHBOARD_URL.txt 2>/dev/null || true
            # shellcheck disable=SC2046,SC2086
            print_link_box_multi $(printf '%s' "$urls")
        fi
    fi

    exec cloudflared tunnel --no-autoupdate run --token "$ZEB_TUNNEL_TOKEN"
fi

# ---------- quick tunnel (fallback, random URL) ----------
echo "[tunnel] starting Cloudflare quick tunnel -> ${target}"

# Clean up cloudflared (and this process group) on stop so a service
# restart doesn't leave an orphaned tunnel behind.
trap 'kill 0 2>/dev/null' TERM INT

last_url=""

# cloudflared prints the assigned URL to its logs. Pipe its output through a
# reader that (a) passes everything through to the container logs and
# (b) reprints the URL in a big box the moment it appears. Keeping the
# pipeline in the foreground lets s6 supervise it and restart on exit.
cloudflared tunnel --no-autoupdate --url "$target" 2>&1 | \
    while IFS= read -r cf_line; do
        printf '%s\n' "$cf_line"
        case "$cf_line" in
            *trycloudflare.com*)
                url=$(printf '%s\n' "$cf_line" \
                    | grep -oE 'https://[a-zA-Z0-9._-]+\.trycloudflare\.com' \
                    | head -n1)
                if [ -n "$url" ] && [ "$url" != "$last_url" ]; then
                    last_url="$url"
                    printf '%s' "$url" > /opt/data/DASHBOARD_URL.txt 2>/dev/null || true
                    print_link_box "$url"
                fi
                ;;
        esac
    done
